# .config/systemd/user/jihou.timer
[Unit]
Description=Jihou Timer

[Timer]
OnCalendar=*-*-* 17:00:00
AccuracySec=1s

[Install]
WantedBy=timers.target

