# .config/systemd/user/jquake-monitoring.timer
[Unit]
Description=JQuake monitoring timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

